# deploy/gunicorn.conf
[program:ohms]
directory=/home/ubuntu/OHMS/back-end
command = /home/ubuntu/OHMS/back-end/.venv/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/OHMS/back-end/app.sock main.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile = /var/log/gunicorn/gunicorn.out.log

[group:ohms]
programs:ohms
